# config/chroot_local-includes/lib/systemd/system/tor@default.service.d/50-resolv-conf-override.conf
[Service]

# Override /etc/resolv.conf so tor can use a clearnet DNS server to
# resolve hostnames used for pluggable transport and proxies.
BindReadOnlyPaths=/etc/resolv-over-clearnet.conf:/etc/resolv.conf
