# config/chroot_local-includes/lib/systemd/system/tor@default.service.d/fix-obfs4proxy.conf
[Service]

# obfs4proxy can't start otherwise
NoNewPrivileges=no
