# config/chroot_local-includes/lib/systemd/system/tor@default.service.d/writable-etc-tor.conf
[Service]

# Tor Launcher needs to make Tor modify its configuration
ProtectSystem=yes
ReadWriteDirectories=-/etc/tor
