# config/chroot_local-includes/etc/xdg/autostart/save-im-environment.desktop.in
[Desktop Entry]
_Name=tails-save-im-environment
_GenericName=save Desktop IM environment
_Comment=save Desktop IM environment so that autostarted iceweasel can use it
Exec=/usr/local/bin/tails-save-im-environment
Terminal=false
Type=Application
Categories=GNOME;X-GNOME-PersonalSettings;
NoDisplay=true
MimeType=application/x-tails-save-im-environment;
