# config/chroot_local-includes/etc/xdg/autostart/security-check.desktop.in
[Desktop Entry]
_Name=tails-security-check
_GenericName=check Tails known security issues
_Comment=check Tails known security issues
Exec=/usr/local/bin/tails-security-check-wrapper
Terminal=false
Type=Application
Categories=GNOME;X-GNOME-PersonalSettings;
NoDisplay=true
MimeType=application/x-tails-security-check;
