# config/chroot_local-includes/etc/xdg/autostart/virt-notify.desktop.in
[Desktop Entry]
_Name=tails-virt-notify-user
_GenericName=warn the user if Tails is running inside a virtual machine
_Comment=warn the user if Tails is running inside a virtual machine
Exec=/usr/local/bin/tails-virt-notify-user
Terminal=false
Type=Application
Categories=GNOME;X-GNOME-PersonalSettings;
NoDisplay=true
MimeType=application/x-tails-virt-notify-user;
