# config/chroot_local-includes/lib/systemd/system/gdm.service.d/failure.conf
[Unit]
OnFailure=tails-gdm-failed-to-start.service
