# config/chroot_local-includes/lib/systemd/system/gdm.service.d/permissions.conf
[Service]

# Tails Greeter writes to the tails.persistence file in there
ExecStartPre=/bin/chgrp Debian-gdm /var/lib/live/config/
ExecStartPre=/bin/chmod g+w /var/lib/live/config/
