# config/chroot_local-includes/lib/systemd/system/gdm.service.d/restart.conf
[Service]

# The GDM unit file has a Restart=always directive, which is good in the
# general case. However, it breaks our emergency shutdown on boot medium
# removal feature, so we disable it.
Restart=no
