# config/chroot_local-includes/etc/NetworkManager/dispatcher.d/00-firewall.sh
#!/bin/sh

# Run whenever an interface gets "up", not otherwise:
if [ "$2" != "up" ]; then
   exit 0
fi

/usr/sbin/ferm /etc/ferm/ferm.conf
