# config/chroot_local-includes/usr/lib/systemd/user/desktop.target
[Unit]
Description=Desktop
Requires=default.target
After=default.target
AllowIsolate=yes
