# config/chroot_local-includes/usr/lib/systemd/user/gnome-early-initialization.target
[Unit]
Description=Bits of GNOME EarlyInitialization managed by systemd
Requires=default.target
After=default.target
AllowIsolate=yes
