# config/chroot_local-includes/etc/network/if-up.d/600tor
#! /bin/sh

if pidof tor > /dev/null; then
    killall -s SIGHUP /usr/sbin/tor
else
    /etc/init.d/tor restart
fi
