# config/chroot_local-includes/lib/systemd/system/tails-additional-software-upgrade.path
[Unit]
Description=Trigger upgrade of Additional Software Packages
Documentation=https://tails.boum.org/contribute/design/persistence/
After=tails-additional-software-install.service
After=tor-has-bootstrapped.service
ConditionFileNotEmpty=/live/persistence/TailsData_unlocked/live-additional-software.conf

[Path]
PathExists=/run/live-additional-software/installed
