# config/chroot_local-includes/lib/systemd/system/tails-allow-external-TailsData-access.service
[Unit]
Description=Allow access to external TailsData partitions
Requires=media-amnesia-TailsData.mount
After=media-amnesia-TailsData.mount

[Service]
ExecStart=setfacl -m user:amnesia:rx /media/amnesia/TailsData

[Install]
WantedBy=media-amnesia-TailsData.mount
