# config/chroot_local-includes/lib/systemd/system/tails-restricted-network-detector.service
[Unit]
Description=Detect restricted networks that may block spoofed MAC addresses

[Service]
Type=simple
ExecStart=/usr/local/lib/tails-restricted-network-detector-wrapper
