# config/chroot_local-includes/lib/systemd/system/tails-sdmem-on-media-removal.service
[Unit]
Description=Wipe memory on live media removal

[Service]
Type=simple
ExecStart=/usr/local/sbin/udev-watchdog-wrapper
After=memlockd.service tails-reconfigure-kexec.service tails-reconfigure-memlockd.service

[Install]
WantedBy=multi-user.target
