# config/chroot_local-includes/lib/systemd/system/tca-portal.socket
[Unit]
Description=TCA portal socket

[Socket]
ListenStream=/run/tca-portal.sock
SocketMode=0600

[Install]
WantedBy=sockets.target

