# config/chroot_local-includes/lib/systemd/system/tca-portal.service
[Unit]
Description=Provides privileged functions for TCA
Requires=tca-portal.socket

[Service]
Type=notify
User=root
ExecStart=/usr/local/lib/tca-portal --systemd-socket
