# wiki/src/todo/htpdate_init_script_Debian_Wheezy_compatibility.mdwn
Please make the htpdate init script compatible with Debian Wheezy. I hacked together a version compatible with Wheezy. If you could test if these changes still work with Tails, that would be awesome!

    #! /bin/sh
    ### BEGIN INIT INFO
    # Provides:          htpdate
    # Required-Start:    mountkernfs $local_fs $remote_fs $syslog networking
    # Required-Stop:     $remote_fs $local_fs $syslog
    # Default-Start:     2 3 4 5
    # Default-Stop:      0 1 6
    # Short-Description: Set time using HTP
    # Description:       Set time using HTP
    ### END INIT INFO

> Sure, we'll have to do that at some point. Thanks for pointing it out.
>
> Please note that while this wiki is fine for tracking tickets, the
> tails-dev ML is better suited for discussing proposed patches.
> Therefore, I'll reply with code review results once the discussion
> is started there.
>
> Bonus points if every proposed atomic change is explained.

[[Wishlist]]
