# config/chroot_local-includes/lib/systemd/system/user@1000.service.d/timeout.conf
[Service]
# Don't let the Upgrader block the reboot after applying an upgrade
TimeoutStopSec=10
